# crontab
* * * * * cd /var/www/project && php artisan schedule:run >> /dev/stdout 2>&1

